# root/etc/systemd/system/zabbix-agent2.service.d/nethserver.conf
[Service]
RuntimeDirectory=zabbix
RuntimeDirectoryMode=0755
